# store-16-master/chef/cookbooks/baseconfig/files/default/rc.local
#!/bin/sh
sleep 10
/usr/local/bin/uwsgi --ini /home/vagrant/project/store16/uwsgi.ini --daemonize /var/log/store16.log